# crontab
1-59/2 * * *   *     cronic ~/bin/bakenpush.sh ~/Documents/diary/ 1
